# s6-rc.d/ollama/run
#!/usr/bin/with-contenv sh
exec ollama serve
